# docker/root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/usr/bin/with-contenv bash

echo -e "\nStarting cron"
/usr/sbin/crond -b -S -l 0 -L /config/cron.log